# etc/supervisord.conf
[inet_http_server]
port = 127.0.0.1:9001
username = yuan
password = 12345


[supervisord]
logfile = data/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = warn
pidfile = data/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
user = lepture  ; you should change this


[program:gunicorn]
command = gunicorn run:app -c etc/gunicorn.py
stdout_logfile = data/gunicorn-stdout.log
redirect_stderr = true

;[program:elasticsearch]
;command=elasticsearch -f etc/elasticsearch.yml
;stdout_logfile = data/elasticsearch-stdout.log
;redirect_stderr = true
